# ircbot.service
[Unit]
Description=IRC Book Downloader Web Interface
After=network.target

[Service]
# --- EDIT THESE TO MATCH YOUR SYSTEM ---
User=YOUR_USER
Group=YOUR_GROUP
WorkingDirectory=/path/to/ircbot
# ----------------------------------------

# Use the uv-managed virtualenv's gunicorn directly
ExecStart=/path/to/ircbot/.venv/bin/gunicorn \
    --worker-class eventlet \
    -w 1 \
    --bind 0.0.0.0:8003 \
    --access-logfile - \
    --error-logfile - \
    app:app

# Environment
Environment="PATH=/path/to/ircbot/.venv/bin:/usr/local/bin:/usr/bin"

# Restart policy
Restart=always
RestartSec=5

# Logging (goes to journalctl)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ircbot

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
